# example/full-node/node.conf
[program:redis]
command=/bin/bash /gryadka/run-redis.sh
autorestart=true
startretries=10

[program:http-proposer]
command=/bin/bash /gryadka/run-http-proposer.sh
autorestart=true
startretries=10