# example/node-acceptor/node.conf
[program:redis]
command=/bin/bash /gryadka/run-redis.sh
autorestart=true
startretries=10