# http-example/node-proposer/node.conf
[program:http-proposer]
command=/bin/bash /gryadka/run-http-proposer.sh
autorestart=true
startretries=10